# services/chat7.service
[Unit]
Description=Chat service 7
Wants=network.target
After=syslog.target network-online.target

[Service]
Type=simple
ExecStart=/opt/chats/chat7
Restart=on-failure
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target
